# scrap/crontab
# START CRON JOB
* * * * * scrapy crawl test_imdb >> /var/log/scrap_top10.log 2>&1
# END CRON JOB